# deploy/nanoclaw.service
[Unit]
Description=NanoClaw Agent Orchestrator
After=network.target docker.service
Requires=docker.service

[Service]
Environment=CONTAINER=1
Environment=PLAYWRIGHT_CHROMIUM_SANDBOX=0
Environment=CHROME_FLAGS=--no-sandbox
Type=simple
User=root
WorkingDirectory=/opt/nanoclaw
# Run delegate patch before every start — survives npm run build
ExecStart=/usr/bin/node --require ./sentry-init.cjs dist/index.js
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3
EnvironmentFile=/opt/nanoclaw/.env
Environment=NODE_ENV=production
# UMask 0000: IPC files created by root are world-writable so containers (uid 1001) can unlink them
UMask=0000
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nanoclaw

[Install]
WantedBy=multi-user.target
